# infra/doldenblick-accounts.service
[Unit]
Description=DoldenBlick Accounts (passwortlose Identität + Onboarding)
Documentation=https://github.com/jihlenburg/hallertau
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=doldenblick
Group=doldenblick
WorkingDirectory=/opt/doldenblick-accounts
Environment=NODE_ENV=production
Environment=PORT=8789
Environment=HOST=127.0.0.1
# Secrets (DATABASE_URL, SESSION_SIGNING_KEY, RP_ID, RP_ORIGIN, SITE_URL, POSTMARK_SERVER_API_TOKEN)
# werden vom Infisical-Sync gerendert (root:root 600). systemd liest sie vor dem Sandboxing.
EnvironmentFile=/etc/doldenblick/doldenblick-accounts.env
ExecStart=/usr/bin/node /opt/doldenblick-accounts/dist/server.js
Restart=always
RestartSec=2

# Härtung — zustandsbehaftet: spricht lokales Postgres (AF_UNIX/AF_INET) + Postmark (AF_INET, HTTPS).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
