# infra/infisical/doldenblick-secrets-sync.timer
[Unit]
Description=DoldenBlick Secrets-Sync alle 10 min (Best-Effort; rs bootet unabhängig vom lokalen File)

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
