# infra/infisical/doldenblick-secrets-sync.service
[Unit]
Description=DoldenBlick: Secrets aus Infisical -> EnvironmentFile syncen (fail-safe, kein Boot-Gate)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/doldenblick-secrets-sync.py apply
# Bewusst KEIN Restart=, keine Abhängigkeit von doldenblick-rs/api — reine Sync-Quelle.
